# campshub360.service
[Unit]
Description=CampsHub360 Django Application
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/ubuntu/campushub-backend-2
EnvironmentFile=/home/ubuntu/campushub-backend-2/.env
Environment=DJANGO_SETTINGS_MODULE=campshub360.production_http
ExecStart=/home/ubuntu/campushub-backend-2/venv/bin/gunicorn --config /home/ubuntu/campushub-backend-2/gunicorn.conf.py --bind 127.0.0.1:8000 campshub360.wsgi:application
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
